# setup/node/sync-data.timer
[Unit]
Description=Sync /data to /mnt daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
